# systemd/msdb-raspberryapp.service
[Unit]
Description=MSDB-RaspberryApp display client
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/msdb-raspberryapp
ExecStart=/opt/msdb-raspberryapp/venv/bin/python -m msdb_raspberryapp
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
