# systemd/msdbr.service
[Unit]
Description=MSDBR display client
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/msdbr
ExecStart=/opt/msdbr/venv/bin/python -m msdbr
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
